# build-support/crontab
@reboot /bin/sh -c '/usr/bin/xbps-reconfigure -a && /usr/bin/rm /var/spool/cron/root'
